# web/99-lan-recovery.rules
ACTION=="add", SUBSYSTEM=="net", KERNEL=="eth0", RUN+="/usr/bin/systemctl restart 
fix-eth.service"
